# supervisord.conf

[program:redis]
command=redis-server

[program:bash]
command=sh ./tsrun.sh

[supervisorctl]
serverurl=http://127.0.0.1:9001
