# supervisord.conf

[program:redis]
command=redis-server

[program:bash]
command=sh ./tsrun.sh
